# deploy/systemd/bumblebee-daily-scan.service
[Unit]
Description=Bumblebee daily deep scan runner

[Service]
Type=oneshot
Environment=BUMBLEBEE_CONFIG=/ABSOLUTE/PATH/TO/REPO/scripts/config.env
ExecStart=/bin/bash /ABSOLUTE/PATH/TO/REPO/scripts/bumblebee-daily-scan.sh

