# deploy/systemd/bumblebee-catalog-watch.service
[Unit]
Description=Bumblebee catalog watcher and scan runner

[Service]
Type=oneshot
Environment=BUMBLEBEE_CONFIG=/ABSOLUTE/PATH/TO/REPO/scripts/config.env
ExecStart=/bin/bash /ABSOLUTE/PATH/TO/REPO/scripts/bumblebee-catalog-watch.sh

