# deploy/systemd/bumblebee-daily-scan.timer
[Unit]
Description=Run Bumblebee daily deep scan

[Timer]
OnCalendar=*-*-* 03:45:00
Persistent=true

[Install]
WantedBy=timers.target
